# release/linux/fpm.before-remove.sh
systemctl stop cerberusd.service
systemctl disable cerberusd.service
