# release/linux/fpm.after-install.sh
systemctl enable cerberusd.service
systemctl start cerberusd.service
